# systems/gate_vm/systemd/gate-broker.service
[Unit]
Description=Gate Lock Broker (Shared Concurrency Substrate)
After=network.target

[Service]
Type=simple
EnvironmentFile=/etc/gate/gate.env
WorkingDirectory=${GATE_REPO_DIR}
ExecStart=/bin/sh -c 'if [ -x /opt/gate/bin/gate ]; then /opt/gate/bin/gate broker --state-dir "${GATE_STATE_DIR}" --host "${GATE_BROKER_HOST}" --port "${GATE_BROKER_PORT}" --lease-ms "${GATE_LEASE_MS}" --max-hold-ms "${GATE_MAX_HOLD_MS}"; else /usr/bin/python3 "${GATE_REPO_DIR}/scripts/gate_broker.py" --state-dir "${GATE_STATE_DIR}" --host "${GATE_BROKER_HOST}" --port "${GATE_BROKER_PORT}" --lease-ms "${GATE_LEASE_MS}" --max-hold-ms "${GATE_MAX_HOLD_MS}"; fi'
Restart=always
User=gate
Group=gate

[Install]
WantedBy=multi-user.target
