# systems/gate_vm/systemd/gate-fuse.service
[Unit]
Description=Gate FUSE Mount (Shared Concurrency Substrate)
After=network.target gate-broker.service
Requires=gate-broker.service

[Service]
Type=simple
EnvironmentFile=/etc/gate/gate.env
WorkingDirectory=${GATE_REPO_DIR}
ExecStartPre=/bin/mkdir -p ${GATE_MOUNT_DIR}
ExecStartPre=/bin/chown gate:gate ${GATE_MOUNT_DIR}
ExecStart=/bin/sh -c 'if [ -x /opt/gate/bin/gate ]; then /opt/gate/bin/gate mount --root "${GATE_REPO_DIR}" --mount "${GATE_MOUNT_DIR}" --broker-host "${GATE_BROKER_HOST}" --broker-port "${GATE_BROKER_PORT}"; else /usr/bin/python3 "${GATE_REPO_DIR}/scripts/gate_mount.py" --root "${GATE_REPO_DIR}" --mount "${GATE_MOUNT_DIR}" --broker-host "${GATE_BROKER_HOST}" --broker-port "${GATE_BROKER_PORT}"; fi'
Restart=always
User=gate
Group=gate

[Install]
WantedBy=multi-user.target
